# ansible/roles/gitea-postgresql/templates/etc/systemd/postgresql.service.j2
{{ ansible_managed | comment }}
.include {{ __postgres_service_file }}
[Service]
Environment=PGDATA={{ __postgres_data_dir }}